# systemd/autostream-updater.service.example
[Unit]
Description=AutoStream central update controller
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=autostream-updater
Group=autostream-updater
ExecStart=/usr/local/bin/autostream-updater run --config /etc/autostream/updater.json
Restart=on-failure
RestartSec=5s
TimeoutStopSec=150s
KillMode=control-group
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectHostname=true
ProtectClock=true
ProtectKernelLogs=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallArchitectures=native
ProtectProc=invisible
ProcSubset=pid
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
ReadWritePaths=/var/lib/autostream-updater

[Install]
WantedBy=multi-user.target
